# build/data/linux/install_files/daemon/upstart/perfect-privacy.conf
# perfect-privacy-vpnd - The daemon for the Perfect Privacy VPN client
#
# Copyright (C) 2014-2015 Perfect Privacy <support@perfect-privacy.com>

description "Perfect Privacy VPN daemon"

start on runlevel [2345]
stop on runlevel [!2345]

expect daemon
respawn
respawn limit 5 10

console output

exec /opt/perfect-privacy/perfect-privacy daemon start
