# build/data/macos/installer/perfect-privacy-service.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<plist version="1.0">
<dict>
    <key>KeepAlive</key>
    <true/>
    <key>inetdCompatibility</key>
    <false/>
    <key>Label</key>
    <string>com.perfect-privacy.service</string>
    <key>ProgramArguments</key>
    <array>
        <string>/Applications/Perfect Privacy.app/Contents/MacOS/perfect-privacy-service</string>
        <string>service</string>
    </array>
    <key>RunAtLoad</key>
    <true/>
    <key>SessionCreate</key>
    <false/>
    <key>UserName</key>
    <string>root</string>
</dict>
</plist>